# smokeping-child/root/etc/services.d/smokeping/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid abc  /usr/bin/smokeping --master-url=$SMOKEPING_MASTER_URL --cache-dir=/var/smokeping --shared-secret=/config/secret.txt --nodaemon
